# aircast/rootfs/etc/s6-overlay/s6-rc.d/esphome-bridge/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: AirCast
# Runs the ESPHome AirPlay Bridge via Shairport-Sync (if enabled)
# ==============================================================================

# Check if ESPHome support is enabled
if bashio::config.false 'esphome_enabled'; then
    bashio::log.info 'ESPHome support is disabled, skipping ESPHome bridge'
    exec sleep infinity
fi

bashio::log.info 'Starting ESPHome AirPlay Bridge (Shairport-Sync)...'

# Basic diagnostics
if bashio::supervisor.ping; then
    bashio::log.info 'Supervisor API reachable'
else
    bashio::log.warning 'Supervisor API NOT reachable (ping failed)'
fi

if bashio::var.has_value "${SUPERVISOR_TOKEN}"; then
    bashio::log.info 'SUPERVISOR_TOKEN present in environment'
else
    bashio::log.warning 'SUPERVISOR_TOKEN missing from environment'
fi

PYTHON_PATH="$(command -v python3 || echo 'python3 not found')"
bashio::log.info "python3 executable: ${PYTHON_PATH}"
bashio::log.info 'Skipping D-Bus and Avahi startup (not required)'

# Run the Shairport-Sync manager for ESPHome devices
bashio::log.info 'Launching shairport-manager.py...'
exec /usr/bin/python3 /usr/bin/shairport-manager.py
